# repro/pkg/fedora/repro.service
[Unit]
Description=reSIProcate repro SIP proxy daemon
After=syslog.target
After=network.target

[Service]
User=repro
Group=repro
Type=forking
ExecStartPre=/bin/mkdir -p /tmp/repro
ExecStartPre=/bin/chown repro:repro /tmp/repro
PIDFile=/tmp/repro/repro.pid
ExecStart=/usr/sbin/repro /etc/repro/repro.config --Daemonize=true --LoggingType=syslog --PidFile=/tmp/repro/repro.pid

[Install]
WantedBy=multi-user.target
